# demos/fdo-serverless-autoupdate/fdo-configs/device0/cfg/var/home/admin/.config/systemd/user/container-httpd-proxy.socket
[Socket]
ListenStream=0.0.0.0:8080
FreeBind=true

[Install]
WantedBy=sockets.target
