# demos/fdo-serverless-autoupdate/fdo-configs/device0/cfg/var/home/admin/.config/systemd/user/pre-pull-container-image.service
[Service]
Type=oneshot
ExecStart=podman pull quay.io/luisarizmendi/2048:prod


[Install]
WantedBy=multi-user.target default.target
