# demos/fdo-serverless-autoupdate/fdo-configs/device0/cfg/etc/systemd/system/applyupdate.service
[Unit]
Description=Apply Update Check

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'if [[ $(rpm-ostree status -v | grep "Staged: yes") ]]; then systemctl --message="Applying OTA update" reboot; else logger "Running latest available update"; fi'

