# demos/fdo-serverless-autoupdate/fdo-configs/device0/cfg/etc/systemd/system/enable-linger.service
[Service]
Type=oneshot
ExecStart=loginctl enable-linger core

[Install]
WantedBy=multi-user.target default.target
